# apps/conf-app/vmconffiles/etc/udev/rules.d/felix.rules
KERNEL=="eth1", SUBSYSTEM=="net", SUBSYSTEMS=="usb", RUN+="/etc/enable-eth.sh"
